# logger/mqtt_service.conf

description "Team254 MQTT Service"

start on filesystem
setuid ubuntu

exec mosquitto -c /home/ubuntu/FRC-2016/logger/mosquitto.conf
